# files/etc/supervisord.d/bucky.ini
[program:bucky]
user=root
command=/usr/bin/bucky /etc/bucky/bucky.conf
environment=PATH="usr/sbin:/usr/bin:/sbin:/bin"
autostart=true
autorestart=true
